# lib/systemd/system/td-agent-bit-reload.path
# Copyright (c) 2019, AT&T Intellectual Property.
# All rights reserved.
#
# SPDX-License-Identifier: GPL-2.0-only
#

[Unit]
Description=td-agent-bit reload on hostname change
Before=td-agent-bit@.service

[Path]
PathChanged=/etc/hostname

[Install]
WantedBy=multi-user.target
