# lib/systemd/system/td-agent-bit-reload.service
[Unit]
Description=td-agent-bit reload

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-reload-or-restart td-agent-bit@*.service
